# systemd/kaplabilling-backend.service
[Unit]
Description=KaplaBilling — FastAPI backend
After=network.target mariadb.service
Requires=mariadb.service

[Service]
Type=simple
User=kaplabilling
Group=kaplabilling
WorkingDirectory=__INSTALL_DIR__/backend
EnvironmentFile=__INSTALL_DIR__/backend/.env
ExecStart=__INSTALL_DIR__/venv/bin/uvicorn main:app --host 127.0.0.1 --port 8000 --workers 2
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=kaplabilling-backend

[Install]
WantedBy=multi-user.target
