# systemd/voxikam-frontend.service
[Unit]
Description=VoxiKam — Next.js frontend
After=network.target voxikam-backend.service

[Service]
Type=simple
User=voxikam
Group=voxikam
WorkingDirectory=__INSTALL_DIR__/frontend
Environment=PORT=3000
Environment=NODE_ENV=production
Environment=HOSTNAME=127.0.0.1
ExecStart=/usr/bin/node __INSTALL_DIR__/frontend/.next/standalone/server.js
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=voxikam-frontend

[Install]
WantedBy=multi-user.target
